# attic/iptables-country-block/etc/network/if-pre-up.d/firewall
#!/bin/sh
/sbin/iptables-restore < /etc/firewall/iptables.rules
