# attic/quirks-no-wifi-powersave [OLD METHOD]/etc/NetworkManager/dispatcher.d/02-wlan-powersave-off
#!/bin/sh
IF=$1
STATUS=$2
if [ "${IF}" = "wls8" ] && [ "${STATUS}" = "up" ]; then
   iwconfig ${IF} power off
   logger "${IF}: turning off powersave mode to prevent constant reconnections"
fi
